# plugins/o9k-roster/systemd/o9k-usage-watcher.service
[Unit]
Description=o9k-roster adaptive subscription usage watcher
After=default.target

[Service]
Type=simple
Restart=always
RestartSec=10
# Symlink after install:
#   install -m 0755 ~/.local/bin/o9k-usage-watcher <path>/o9k-usage-watcher
ExecStart=%h/.local/bin/o9k-usage-watcher

[Install]
WantedBy=default.target
