# plugins/o9k-roster/systemd/o9k-resume.service
[Unit]
Description=o9k agentless cross-CLI run resume after boot
After=default.target

[Service]
Type=oneshot
# Prefer a stable path after install. Until o9k-init rewrites this, symlink:
#   mkdir -p ~/.local/bin && ln -sf /path/to/plugins/o9k-roster/scripts/runs.mjs ~/.local/bin/o9k-runs
ExecStart=%h/.local/bin/o9k-runs resume

[Install]
WantedBy=default.target
