# plugins/o9k-roster/systemd/o9k-usage-watcher.service
[Unit]
Description=o9k-roster adaptive subscription usage watcher
After=default.target

[Service]
Type=simple
Restart=always
RestartSec=10
# Install wrapper (never ExecStart a bare .mjs — ENOEXEC without shebang):
#   install -m 0755 <repo>/plugins/o9k-roster/scripts/o9k-usage-watcher.sh ~/.local/bin/o9k-usage-watcher
# Or symlink the repo script if it stays on a fixed path.
ExecStart=%h/.local/bin/o9k-usage-watcher

[Install]
WantedBy=default.target
