# plugins/o9k-roster/systemd/o9k-usage-watcher.service
[Unit]
Description=o9k-roster adaptive subscription usage watcher
After=default.target

[Service]
Type=simple
Restart=always
RestartSec=10
# Install (pick one):
#   ln -sf <repo>/plugins/o9k-roster/scripts/o9k-usage-watcher.sh ~/.local/bin/o9k-usage-watcher
#   # copy + env (foreign install path):
#   install -m 0755 <repo>/plugins/o9k-roster/scripts/o9k-usage-watcher.sh ~/.local/bin/o9k-usage-watcher
#   mkdir -p ~/.config/systemd/user/o9k-usage-watcher.service.d
#   printf '[Service]\nEnvironment=O9K_ROSTER_SCRIPTS=<repo>/plugins/o9k-roster/scripts\n' \
#     > ~/.config/systemd/user/o9k-usage-watcher.service.d/path.conf
#   systemctl --user daemon-reload
ExecStart=%h/.local/bin/o9k-usage-watcher

[Install]
WantedBy=default.target
